# runit/log/run
#!/bin/sh

mkdir -p /var/log/publisher
exec svlogd -tt /var/log/publisher


